# src/systemd/arch-update.service
[Unit]
Description=Run arch-update auto check
After=network-online.target

[Service]
Type=oneshot
ExecStart=arch-update --check
