# systemd/arch-check-update-notifier.service
[Unit]
Description=Arch Check Update Notifier
After=network-online.target

[Service]
Type=oneshot
ExecStart=arch-check-update-notifier