# systemd/arch-update.service
[Unit]
Description=Arch Update auto check
After=network-online.target

[Service]
Type=oneshot
ExecStart=arch-update --check
